# kiosk-monitor.service
[Unit]
Description=Kiosk Monitor - Restart X on crash
After=getty@tty1.service

[Service]
Type=simple
ExecStart=/home/realo/kiosk-monitor.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
